# config/systemd/user/yubikey-glow.service
[Unit]
Description=yubikey-glow rainbow screen-edge glow (OpenPGP/YubiKey touch indicator)
# Reads yubikey-touch-detector's socket; the daemon retries the connection on
# its own, so this is a soft ordering hint, not a hard requirement.
Wants=yubikey-touch-detector.socket
After=yubikey-touch-detector.socket

[Service]
Type=simple
# Needs WAYLAND_DISPLAY/XDG_RUNTIME_DIR in the manager environment. sway's
# `systemctl --user import-environment` (config line 2) supplies them, but sway
# fires its exec lines concurrently, so that import may not have landed yet on
# the first start — the script then exits nonzero on the missing display and
# Restart=on-failure brings it back up once the env is in place.
ExecStart=%h/git/kannar/nunux/bin/yubikey-glow
Restart=on-failure
RestartSec=2

# Not enabled on a target: graphical-session.target isn't activated on this box,
# so sway starts the unit explicitly. No [Install] section by design.
